# psad-init.gentoo
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$

depend() {
	need net logger
	use dns
}

PIDFILE=/var/run/psad/psad.pid

start() {
	ebegin "Starting psad daemons"
	/usr/sbin/psad &>/dev/null
	eend $?
}

stop() {
	ebegin "Stopping psad daemons"
	/usr/sbin/psad --Kill &>/dev/null
	eend $?
}
